# v2/99-usb-stick.rules
ACTION=="add",ATTRS{idVendor}=="13fe",ATTRS{idProduct}=="6300",SYMLINK+="datadump",RUN+="/home/ss/copyfiles"
